# nvr-demo.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/exec-demo.sh

[Install]
WantedBy=multi-user.target
