# services/be-go.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/go-exec.sh

[Install]
WantedBy=multi-user.target
