# services/be-go.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=$NVRPATH/scripts/go-exec.sh

[Install]
WantedBy=multi-user.target
