# services/be-nodejs.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/nodejs-exec.sh

[Install]
WantedBy=multi-user.target
