# services/be-python.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/python-exec.sh

[Install]
WantedBy=multi-user.target
