# services/fe-reactjs.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/usr/bin/fe-exec.sh

[Install]
WantedBy=multi-user.target
