# services/fe-reactjs.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=$NVRPATH/scripts/fe-exec.sh

[Install]
WantedBy=multi-user.target
